# prod/viewlinc-telegram.service
[Unit]
Description=ViewLinc Telegram Bot with Time-Series Plotting
After=network.target

[Service]
Type=simple
User=river
WorkingDirectory=/home/river/Desktop/monitoring
ExecStart=/usr/bin/Rscript /home/river/Desktop/monitoring/telegram_bot_linux.R
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
